# services/little-timmy.service
[Unit]
Description=Little Timmy Voice Assistant Orchestrator (port 8893)
After=network-online.target postgresql.service whisper-server.service llama-3b-server.service qwen36-server.service ollama.service
Wants=network-online.target
Requires=postgresql.service

[Service]
Type=simple
User=gearscodeandfire
Group=gearscodeandfire
WorkingDirectory=/home/gearscodeandfire/little_timmy
ExecStart=/home/gearscodeandfire/little_timmy/.venv/bin/python main.py
Restart=on-failure
RestartSec=5
KillMode=process
TimeoutStopSec=10
Environment=PYTHONUNBUFFERED=1
# Route audio through the user's PipeWire/Pulse session so Piper's 22050 Hz
# output is resampled instead of failing on a raw ALSA device. UID 1000 =
# gearscodeandfire. Without these, sd.play() opens the system default ALSA
# device which only accepts 44100/48000 -> PaErrorCode -9997.
Environment=XDG_RUNTIME_DIR=/run/user/1000
Environment=PULSE_SERVER=unix:/run/user/1000/pulse/native
Environment=TIMMY_VISION_URL=http://localhost:8084

[Install]
WantedBy=multi-user.target
