# services/qwen3-4b-classifier.service
[Unit]
Description=llama.cpp Server (Qwen3-4B-Q4_K_M, port 8092, first-pass tool-call/intent classifier)
Documentation=https://github.com/ggerganov/llama.cpp
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=gearscodeandfire
Group=gearscodeandfire
ExecStart=/usr/bin/stdbuf -oL /home/gearscodeandfire/llama.cpp/build-vulkan-head/bin/llama-server \
    -m /home/gearscodeandfire/models/Qwen3-4B-Q4_K_M.gguf \
    --host 0.0.0.0 \
    --port 8092 \
    -ngl 99 \
    -c 4096 \
    -fa on -ctk q8_0 -ctv q8_0 \
    -cb -np 1 \
    --reasoning-format deepseek --jinja
Restart=on-failure
RestartSec=10
StandardOutput=append:/home/gearscodeandfire/demerzel/logs/qwen3-4b-classifier.log
StandardError=append:/home/gearscodeandfire/demerzel/logs/qwen3-4b-classifier.log
KillMode=process
TimeoutStopSec=10

[Install]
WantedBy=multi-user.target
